# files/deb/prerm
#!/bin/sh
# prerm script for nsclient
#

if [ -x "/etc/init.d/nsclient" ] || [ -e "/etc/init/nsclient.conf" ]; then
        invoke-rc.d nsclient stop || exit $?
fi
